# tools/osbuilder/rootfs-builder/datadog-files/etc/systemd/system/system-probe.service.d/50-apparmor.conf
[Service]
AppArmorProfile=datadog-system-probe
