# tools/osbuilder/rootfs-builder/datadog-files/etc/systemd/system/kata-agent.service.d/50-apparmor.conf
[Service]
AppArmorProfile=usr.bin.kata-agent
